# deploy/lingbase.service
[Unit]
Description=Lingbase Edge LLM Inference Service
Documentation=https://github.com/your-org/lingbase
After=network.target

[Service]
Type=simple
User=firefly
WorkingDirectory=/home/firefly/lingbase/lingbase-0.1.0-aarch64
ExecStart=/home/firefly/lingbase/lingbase-0.1.0-aarch64/lingbase
Restart=always
RestartSec=10

# Environment
Environment=RUST_LOG=info

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lingbase

# Resource limits
LimitNOFILE=65536
MemoryMax=16G

[Install]
WantedBy=multi-user.target